# deploy/nextg3n-mcp.service
[Unit]
Description=NextG3N MCP Manager Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/nextg3n
ExecStart=/usr/bin/python -m mcp.mcp_manager
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nextg3n-mcp
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target